# platform/proxmox/nodes/pve-guide/etc/systemd/system/pve-guests.service.d/override.conf
[Unit]
Requires=nvidia-devnodes.service
After=nvidia-devnodes.service
